# services/collector/crontabs/daily
# Daily collectors
0 4 * * * TASK=camera npx tsx index.ts 2>&1 && touch /tmp/collector-healthy
0 6 * * * TASK=charger npx tsx index.ts 2>&1
30 22 * * * TASK=insights npx tsx index.ts 2>&1
