# services/collector/crontabs/daily
# Daily collectors
# Staggered schedule — no two heavy tasks overlap.
# Heaviest: charger (768MB heap), mobilitydata-sync (ZIP parsing, ~512MB peak)

30 0 * * * TASK=daily-stats npx tsx index.ts 2>&1

# CNMC provincial fuel price history — moved to weekly (CNMC publishes weekly aggregates)
# See crontabs/weekly for the Sunday 02:00 entry.

# MobilityData GTFS catalog sync — moved to weekly (catalog updates infrequently)
# See crontabs/weekly for the Sunday 02:30 entry.

# INE transport statistics — INE publishes quarterly (modal split, annual
# series). Daily polling re-fetched the same data ~28 of 30 days. Now Sunday
# 03:30 weekly; still well inside INE's quarterly cadence with headroom.
30 3 * * 0 TASK=ine-stats npx tsx index.ts 2>&1

0 4 * * * TASK=camera npx tsx index.ts 2>&1 && touch /tmp/last-run

# AENA airport passenger stats — Eurostat AVIA_PAOA publishes annually with
# a ~12-month lag. Daily polling was pure noise. Now monthly on the 1st at
# 04:30 (overkill for annual data but a useful low-bound liveness signal).
30 4 1 * * TASK=aena-stats npx tsx index.ts 2>&1

0 5 * * * SYNC_MODE=incremental TASK=typesense-sync npx tsx index.ts 2>&1
0 6 * * * TASK=charger NODE_OPTIONS=--max-old-space-size=768 npx tsx index.ts 2>&1

# AEMET daily climate update (after station data stabilizes for previous day)
0 8 * * * TASK=aemet-historical npx tsx index.ts 2>&1

# AEMET 7-day forecast per municipio (4x/day, staggered to avoid peaks) — T3.3
0 3,9,15,21 * * * TASK=aemet-forecast npx tsx index.ts 2>&1

# CAMS AQ 5-day forecast per gridpoint (2x/day) — T3.4
30 6,18 * * * TASK=cams-aq npx tsx index.ts 2>&1

30 22 * * * TASK=insights npx tsx index.ts 2>&1

# Nightly cleanup — rolling-window realtime tables (TransitVehiclePosition 48h)
0 23 * * * TASK=cleanup-realtime npx tsx index.ts 2>&1

# NOTE: accident-microdata — when added, requires NODE_OPTIONS=--max-old-space-size=512
# (large CSV parsing, ~400MB peak). Example:
# 0 1 * * 0 TASK=accident-microdata NODE_OPTIONS=--max-old-space-size=512 npx tsx index.ts 2>&1
