# services/collector/crontabs/frequent
# Frequent collectors — hourly weather, 6h maritime/portugal
# Weather hourly for storm freshness (AEMET alerts change within minutes)
0 * * * * TASK=weather npx tsx index.ts 2>&1 && touch /tmp/last-run
0 */6 * * * TASK=maritime-forecast npx tsx index.ts 2>&1
0 */6 * * * TASK=portugal-weather npx tsx index.ts 2>&1
# NOTE: air-quality runs in the realtime crontab (hourly at :00) — do NOT add here

# Voyage detector — groups VesselPosition rows into PortCalls & Voyages (hourly)
# Idempotent upsert; requires AIS positions to be flowing (ais-stream container)
# Lockfile guard: skip if a previous run is still in progress (takes 2-4h for 55k vessels)
0 * * * * flock -n /tmp/voyage-detector.lock -c 'TASK=voyage-detector npx tsx index.ts 2>&1' || echo "[voyage-detector] skipped — previous run still in progress"

# SASEMAR: historical archive (2019-2024), one-shot only — do NOT re-enable on cron
# SASEMAR maritime emergencies — historical ArcGIS upsert (every 15 min)
# Source: opendata-sasemar.hub.arcgis.com (2019-2024 historical records)
# Idempotent: upserts by external ID; safe to run frequently
#*/15 * * * * TASK=sasemar npx tsx index.ts 2>&1
