# services/collector/crontabs/frequent
# Frequent collectors — every 6 hours
# Weather data from AEMET, IPMA
0 */6 * * * TASK=weather npx tsx index.ts 2>&1 && touch /tmp/collector-healthy
0 */6 * * * TASK=maritime-forecast npx tsx index.ts 2>&1
0 */6 * * * TASK=portugal-weather npx tsx index.ts 2>&1
