# services/collector/crontabs/realtime
# Realtime collectors — every 2-5 minutes
# Supercronic runs each line as a separate process.
# Peak concurrency: ~11 tasks at :00/:10/:20/:30/:40/:50 marks
# Budget: ~80MB per Node process × 11 concurrent = ~880MB peak → 1536m limit

# DGT + regional traffic feeds (every 2 min for the multi-source `incident` task,
# which polls 9 regional + DGT NAP DATEX II endpoints; tighter than DGT's own
# 5-min publish cadence so regional feeds drive resolution).
*/2 * * * * TASK=incident npx tsx index.ts 2>&1 && touch /tmp/last-run
# v16 — DGT NAP V16 beacon feed. Source publishes every 5 min, so polling
# at */5 (not */2) cuts our DGT NAP load in half — `incident` already hits
# the same DGT URL every 2 min, so the every-2-min v16 was pure waste.
*/5 * * * * TASK=v16 npx tsx index.ts 2>&1

# DGT infrastructure (every 5 min)
*/5 * * * * TASK=panel npx tsx index.ts 2>&1
*/5 * * * * TASK=detector npx tsx index.ts 2>&1

# Madrid intensity sensors (every 5 min, ~6K measurement points)
*/5 * * * * TASK=intensity npx tsx index.ts 2>&1

# City traffic sensors — Barcelona/Valencia/Zaragoza (every 5 min)
*/5 * * * * TASK=city-traffic npx tsx index.ts 2>&1

# Andorra cross-border (every 5 min)
*/5 * * * * TASK=andorra npx tsx index.ts 2>&1

# Renfe GTFS-RT alerts + fleet positions (every 2 min)
*/2 * * * * TASK=renfe-alerts npx tsx index.ts 2>&1
*/2 * * * * TASK=renfe-ld-realtime npx tsx index.ts 2>&1
*/2 * * * * TASK=renfe-positions npx tsx index.ts 2>&1

# Urban transit GTFS-RT (EMT Madrid, TMB Barcelona) — every 1 min
# Skips individual feeds if API keys missing; no-op when all skip
* * * * * TASK=transit-realtime npx tsx index.ts 2>&1

# OpenSky aircraft positions — anonymous limit 400 req/day (~1 per 3.6 min), authenticated 100 req/min
# Inner loop runs multiple polls if authenticated (uses the 10s minimum); single poll otherwise.
# */4 is safe for anonymous mode (360 req/day < 400 limit); authenticated mode polls internally every 30s.
*/4 * * * * TASK=opensky npx tsx index.ts 2>&1

# MITECO air quality — 506 stations, ICA index (hourly)
0 * * * * TASK=air-quality npx tsx index.ts 2>&1

# EUMETSAT radar composite (15-min, RainViewer fallback if credentials missing) — T3.2
*/15 * * * * TASK=eumetsat-radar npx tsx index.ts 2>&1

# Alert matching engine (every 5 min) — ROAD/TRAIN/FLIGHT × Push/Email/Telegram — T4.8g
*/5 * * * * TASK=alert-matcher npx tsx index.ts 2>&1

# Social broadcast (every 5 min, dedup via Redis 24h) — T4.9d
*/5 * * * * TASK=social-broadcast npx tsx index.ts 2>&1

# Health check failsafe — checks all tables against SLAs, triggers recovery for stale data
*/30 * * * * TASK=health-check npx tsx index.ts 2>&1
