# services/collector/crontabs/weekly
# Weekly collectors — Sunday early morning
# Staggered 30-min intervals. Heaviest: transit-gtfs (145 auto-discovered feeds, Xunta 141MB)

0 3 * * 0 TASK=radar npx tsx index.ts 2>&1 && touch /tmp/last-run
30 3 * * 0 TASK=speedlimit npx tsx index.ts 2>&1
0 4 * * 0 TASK=risk-zones npx tsx index.ts 2>&1
30 4 * * 0 TASK=zbe npx tsx index.ts 2>&1

# Renfe GTFS static (weekly — timetables)
0 5 * * 0 TASK=renfe-gtfs NODE_OPTIONS=--max-old-space-size=768 npx tsx index.ts 2>&1

# Ferry GTFS — Fred. Olsen, Baleària, Vizcaya (3 small feeds)
30 5 * * 0 TASK=ferry-gtfs npx tsx index.ts 2>&1

# Transit GTFS — auto-discovered 145 ES feeds from MobilityData catalog
# Needs 4GB heap for Xunta de Galicia (141MB, 156K trips, 27K stops)
0 6 * * 0 TASK=transit-gtfs NODE_OPTIONS=--max-old-space-size=4096 npx tsx index.ts 2>&1

# DGT extras — roadworks/connected cones (placeholder, probes NAP endpoints)
30 6 * * 0 TASK=dgt-extras npx tsx index.ts 2>&1

# IMD (Intensidad Media Diaria) — Ministry of Transport ArcGIS REST API
# ~14K counting stations + ~14K road segments for years 2017-2019 + 2022.
# Data changes rarely so weekly is more than enough; upsert is idempotent.
# Monthly was specified in docs but the cron was never wired, which left
# TrafficFlow at 0 rows in prod. Restoring as weekly on Sunday night.
0 7 * * 0 TASK=imd NODE_OPTIONS=--max-old-space-size=1024 npx tsx index.ts 2>&1

# Weekly digest newsletter (Monday 09:00) — T4.9c
0 9 * * 1 TASK=weekly-digest npx tsx index.ts 2>&1

# Monthly state-of-traffic PDF report (1st of month 03:00) — T4.9e
0 3 1 * * TASK=monthly-report NODE_OPTIONS=--max-old-space-size=1024 npx tsx index.ts 2>&1
