# services/otp/crontab
# OTP2 graph rebuild — every Monday at 03:00 UTC
# Rebuilds from current OSM PBF + GTFS feeds; hot-swaps graph on success.
# After rebuild, OTP container is restarted to load the new graph.obj.
#
# Log: /var/log/otp-build.log

0 3 * * 1 root /opt/trafico/services/otp/build-graph.sh >> /var/log/otp-build.log 2>&1
